# systemd/voicetypist.service
[Unit]
Description=VoiceTypist Linux Dictation
After=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/opt/VoiceTypist
ExecStart=%h/opt/VoiceTypist/venv/bin/python %h/opt/VoiceTypist/voicetypist_linux.py
Environment=PATH=%h/opt/VoiceTypist/venv/bin:/usr/bin:/bin
PassEnvironment=DISPLAY WAYLAND_DISPLAY XAUTHORITY XDG_CURRENT_DESKTOP XDG_SESSION_TYPE DBUS_SESSION_BUS_ADDRESS GOOGLE_API_KEY
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
